# rootfs/etc/s6-overlay/s6-rc.d/lemonade-ha/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Lemonade HA add-on
# Installs the custom integration into HA config dir, then sleeps.
# All model/connection configuration is done through the HA integration UI.
# ==============================================================================

# ── Find HA config directory ──────────────────────────────────────────────────
# Current HA OS mounts it at /homeassistant; older installs use /config.
if [ -d "/homeassistant" ]; then
    HA_CONFIG="/homeassistant"
else
    HA_CONFIG="/config"
fi

COMPONENT_SRC="/usr/src/custom_components/lemonade_ha"
COMPONENT_DEST="${HA_CONFIG}/custom_components/lemonade_ha"

mkdir -p "${HA_CONFIG}/custom_components"

if [ ! -d "${COMPONENT_DEST}" ] || ! diff -rq "${COMPONENT_SRC}" "${COMPONENT_DEST}" > /dev/null 2>&1; then
    bashio::log.info "Installing/updating Lemonade HA custom integration ..."
    rm -rf "${COMPONENT_DEST}"
    cp -r "${COMPONENT_SRC}" "${COMPONENT_DEST}"
    bashio::log.info "Custom integration installed to: ${COMPONENT_DEST}"
    bashio::log.warning "*** Restart Home Assistant (not the add-on) to load the integration. ***"
    bashio::log.warning "Then: Settings > Devices & Services > Add Integration > search 'Lemonade'."
else
    bashio::log.info "Lemonade HA custom integration is already up to date at: ${COMPONENT_DEST}"
fi

bashio::log.info "Add-on running. Integration is managed by Home Assistant."

# Sleep forever — s6 keeps this process alive; the actual work is done by the
# integration inside HA core, not inside this container.
exec sleep infinity
